# .config/systemd/user/keybase-kbfs.service
[Unit]
Description=Keybase KBFS
After=network.target

[Service]
ExecStart=/usr/bin/kbfsfuse -debug /keybase

[Install]
WantedBy=default.target
