# .config/systemd/user/urxvtd.service
[Unit]
Description=URxvt Daemon Service

[Service]
ExecStart=/usr/bin/urxvtd -q -o

[Install]
WantedBy=default.target
