# nolink/systemd/user/urxvtd.service
[Unit]
Description=URxvt Daemon Service

[Service]
ExecStart=/usr/bin/urxvtd -q -o
Restart=always

[Install]
WantedBy=default.target
