# nolink/.config/systemd/user/keybase.service
[Unit]
Description=Keybase
After=network.target

[Service]
ExecStart=/usr/bin/keybase -d service

[Install]
WantedBy=default.target

