# kitty/.local/share/systemd/user/ysgrifennwr.service
[Unit]
Description=Change kitty colorscheme to light or dark according to the time of the day

[Service]
ExecStart=/usr/bin/zsh -c 'if\x20\x5b\x5b\x20\x24\x28\x28\x20\x28\x28\x24\x28date\x20\x2b\x25\x25\x2dH\x25\x25M\x29\x20\x2b\x201600\x29\x20\x25\x25\x202400\x29\x20\x29\x29\x20\x2dlt\x20830\x20\x5d\x5d\x3b\x20then\x3b\x20kitty\x20\x40\x20\x2d\x2dto\x20unix:\x40mykitty\x20set\x2dcolors\x20\x2da\x20\x2dc\x20foreground\x3d\x22\x23424242\x22\x20background\x3d\x22\x23f9f8f4\x22\x20color0\x3d\x22\x23f9f8f4\x22\x20color1\x3d\x22\x23e32791\x22\x20color2\x3d\x22\x23488432\x22\x20color3\x3d\x22\x23a25d0e\x22\x20color4\x3d\x22\x232c65b5\x22\x20color5\x3d\x22\x23B062A7\x22\x20color6\x3d\x22\x2327BBBE\x22\x20color7\x3d\x22\x23999999\x22\x20color8\x3d\x22\x23B8B8B8\x22\x20color9\x3d\x22\x239F1B66\x22\x20color10\x3d\x22\x23325D23\x22\x20color11\x3d\x22\x2371410A\x22\x20color12\x3d\x22\x231F477F\x22\x20color13\x3d\x22\x237B4474\x22\x20color14\x3d\x22\x231B8486\x22\x20color15\x3d\x22\x23424242\x22\x20cursor\x3d\x22\x23FC9520\x22\x3b\x20else\x3b\x20kitty\x20\x40\x20\x2d\x2dto\x20unix:\x40mykitty\x20set\x2dcolors\x20\x2da\x20\x2dc\x20foreground\x3d\x22\x23e5e6e6\x22\x20background\x3d\x22\x231e1e1e\x22\x20color0\x3d\x22\x231e1e1e\x22\x20color8\x3d\x22\x23515151\x22\x20color1\x3d\x22\x23e32791\x22\x20color9\x3d\x22\x23e566ad\x22\x20color2\x3d\x22\x2330c798\x22\x20color10\x3d\x22\x236cd1b2\x22\x20color3\x3d\x22\x23e3c472\x22\x20color11\x3d\x22\x23e4cf98\x22\x20color4\x3d\x22\x236769e6\x22\x20color12\x3d\x22\x239091e6\x22\x20color5\x3d\x22\x23e59fdf\x22\x20color13\x3d\x22\x23e5b6e1\x22\x20color6\x3d\x22\x2381d8d0\x22\x20color14\x3d\x22\x23a2dcd7\x22\x20color7\x3d\x22\x23969696\x22\x20color15\x3d\x22\x23e5e6e6\x22\x20cursor\x3d\x22\x2320bbfc\x22\x3b\x20fi'
