# systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -a xha -J %I $TERM
TTYVTDisallocate=no